# ops/systemd/atel-180-3300.service.d/hook.conf
[Service]
Environment=ATEL_HOOK_CONCURRENCY=1
