# ops/systemd/atel-auto-requester-3921.service.d/atel-openclaw.conf
[Service]
Environment=ATEL_AGENT_CMD=/usr/local/bin/atel-openclaw-agent.sh
Environment=ATEL_OPENCLAW_AGENT_NAME=main
Environment=ATEL_HOOK_CONCURRENCY=4
